# systemd/memory-ingest.timer
[Unit]
Description=Memory Engine Auto-Ingest Timer (every 5 minutes)

[Timer]
OnBootSec=60
OnUnitActiveSec=300
AccuracySec=30

[Install]
WantedBy=timers.target
